# mongo.sv.conf
[program:mongo]
command=/bin/bash -c "/mongo-entrypoint.sh && mongod --smallfiles"
redirect_stderr=true
