# mongo.sv.conf
[program:mongo]
command=/bin/bash -c "mongod --smallfiles"
redirect_stderr=true
